# agent/templates/pgpool2-exporter.service
[Unit]
Description=pgpool2_exporter (Zabbix observability)
After=network.target
Wants=network.target

[Service]
Type=simple
# Credenciais/host via env file (modo 600, gerado pelo install.sh) - nunca no git.
# Vars usadas pelo pgpool2_exporter: POSTGRES_USERNAME, POSTGRES_PASSWORD,
# PGPOOL_SERVICE, PGPOOL_SERVICE_PORT.
EnvironmentFile=/etc/pgpool2-exporter/pgpool2_exporter.env
ExecStart=/usr/local/bin/pgpool2_exporter --web.listen-address=127.0.0.1:9719
Restart=always
RestartSec=5
User=root
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
